# tamara/services/home.service
[Unit]
Description="Tamara home"
After="network.target"

[Service]
ExecStart="/usr/bin/local/tamara/home.py"

[Install]
WantedBy=multi-user.target
